# system_files/usr/lib/systemd/user/wayvnc.service
[Unit]
Description=WayVNC Remote Desktop Server (localhost-only)
Documentation=man:wayvnc(1)
# Si aggancia alla graphical session: parte con niri, muore con niri
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/wayvnc --config /usr/share/gdots/system/wayvnc/wayvnc.conf
Restart=on-failure
RestartSec=5
# Log strutturato su journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wayvnc

[Install]
WantedBy=graphical-session.target
